# systemd/cc-hub-updater.service
[Unit]
After=docker.service

[Service]
ExecStart=/opt/client-connector-hub/updater.sh
Restart=always

[Install]
WantedBy=default.target
